# bloomberg-bot.service
[Unit]
Description=Bloomberg Alert Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bloomberg-bot
Environment=PATH=/home/ubuntu/bloomberg-bot/venv/bin
ExecStart=/home/ubuntu/bloomberg-bot/venv/bin/python bloomberg_simple.py
Restart=always
RestartSec=10

# Environment variables (alternative to .env file)
Environment=TELEGRAM_TOKEN=8122220616:AAFI8xFd2O0X1UiJWBWvO8j5-pgeysLCbpc
Environment=GEMINI_API_KEY=AIzaSyAPYV-zu5bNVFcTeElaXYWokZEk_wlAWms
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
